# yamix.service
[Unit]
Description=Yamix Web UI
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/yamix
ExecStart=/opt/yamix/bin/web.py

